# .config/systemd/user/backup-daily.timer
[Unit]
Description=Daily backup of $h/df 

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
